# nextcloud/build/nextcloud-fpm/supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/var/run/supervisord/supervisord.pid
logfile=/var/run/supervisord/supervisord.log

[program:php-fpm]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=php-fpm

[program:cron]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/cron.sh
